# anarchistos-configs/data/etc/skel/.config/autostart/Bitmask.desktop
[Desktop Entry]
Name=Bitmask
Type=Application
Exec=/usr/bin/bitmask
Path=/home/areyouloco
Terminal=false
